# supervisord-check_weed.ini
[program:check_weed]
directory=/sunlight/shell/monitor
command=/usr/bin/sh check_weed.sh
autostart=true
startsecs=10
autorestart=true
startretries=30
user=root
redirect_stderr=true
stopasgroup=true
killasgroup=true
loglevel=info
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisord/log/check_weed.log
